# packaging/app-whitelisting.cron
00 23 * * * root /usr/sbin/update-fingerprints >/dev/null 2>&1

